# packaging/proxmox-license-proxy.service
[Unit]
Description=Proxmox License Proxy
Documentation=https://github.com/TheMeinerLP/proxmox-license-proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pmox
Group=pmox
ExecStart=/usr/bin/proxmox-license-proxy serve --config /etc/pmox/config.yaml
Restart=on-failure
RestartSec=5s

# Bind the default privileged port (:443) as a non-root user.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Hardening / sandboxing.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
# AF_NETLINK is required so the binary can enumerate local interface addresses
# (Go's net.Interfaces), which the mDNS advertiser needs - without it startup
# logs "mDNS advertise failed: Could not determine host IP addresses".
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
# Group-writable (0007) so a root-run CLI and the pmox service can share the
# registry in /var/lib/pmox (a setgid, group-pmox directory); the app also pins
# the registry files to 0660 explicitly.
UMask=0007

# Writable state: the server persists registry.json (auto-registered hosts) here.
StateDirectory=pmox
StateDirectoryMode=2770
ReadWritePaths=/var/lib/pmox
ReadOnlyPaths=/etc/pmox

[Install]
WantedBy=multi-user.target
